# ansible/playbook/roles/uwsgi/templates/upstart.uwsgi.conf.j2
# simple uWSGI upstart job

start on runlevel [2345]
stop on runlevel [06]
respawn

exec {{ project.virtualenv }}/bin/uwsgi --ini {{project.root}}/uwsgi.ini
